# playbook-base/files/runit/sv/nftables/finish
#!/bin/sh
set -e
/usr/sbin/nft flush ruleset 2>/dev/null || true
