# playbook-base/files/runit/sv/nftables/run
#!/bin/sh
# Stateless firewall loader: apply the ruleset once with `nft -f` (atomic), then
# `sleep infinity` so runsv keeps this service supervised — the rules live in the
# kernel regardless of this process. (`pause` is not a real command; runit's
# load-and-hold idiom is a blocking sleep.) On stop runsv runs ./finish (flush).
exec 2>&1
[ -r /etc/nftables.conf ] && /usr/sbin/nft -f /etc/nftables.conf
exec sleep infinity
