# playbook-base/files/runit/sv/iptables/finish
#!/bin/sh
# Reset to default-accept on service stop so a sysadmin yanking the
# firewall doesn't lock themselves out of a half-configured box.
set -e
for t in filter nat mangle raw; do
  /sbin/iptables  -t "$t" -F 2>/dev/null || true
  /sbin/iptables  -t "$t" -X 2>/dev/null || true
  /sbin/ip6tables -t "$t" -F 2>/dev/null || true
  /sbin/ip6tables -t "$t" -X 2>/dev/null || true
done
/sbin/iptables  -P INPUT   ACCEPT
/sbin/iptables  -P OUTPUT  ACCEPT
/sbin/iptables  -P FORWARD ACCEPT
/sbin/ip6tables -P INPUT   ACCEPT
/sbin/ip6tables -P OUTPUT  ACCEPT
/sbin/ip6tables -P FORWARD ACCEPT
