# playbook-base/files/runit/sv/nix-daemon/run
#!/bin/sh
# Multi-user nix-daemon. sys-apps/nix on Gentoo lays the daemon under the
# /nix profile tree (not /usr/bin/), and it needs `--daemon` to run in
# foreground (runsv requires foreground processes). The build users
# (nixbld1..N) and /nix/store ownership are set up by the sys-apps/nix
# ebuild's post_install hook — this run script assumes that already
# happened. See https://wiki.gentoo.org/wiki/User:Alxhr0/Nix_on_openrc
exec 2>&1
[ -r conf ] && . ./conf
exec /nix/var/nix/profiles/default/bin/nix-daemon --daemon ${OPTS}
