# playbook-base/files/runit/sv/udevd/run
#!/bin/sh
exec 2>&1
# rc/1.02-udev.sh starts udev as a daemon during stage 1. That instance
# holds /run/systemd/udev/control, so a fresh foreground udevd exec'd
# here would fail "already started by another instance" and runsv would
# crash-loop the service. Tell the stage-1 daemon to exit first, then
# exec our own (which becomes runsv's supervised child).
udevadm control --exit 2>/dev/null || true
# Pick whichever udev daemon the system has — sys-apps/systemd-utils[udev]
# installs systemd-udevd; sys-fs/eudev installs /sbin/udevd.
for _u in /usr/lib/systemd/systemd-udevd /lib/systemd/systemd-udevd /sbin/udevd; do
    [ -x "$_u" ] && exec "$_u"
done
echo "no udevd binary found" >&2
exit 1
